# .service
[Unit]
Description=Laravel queue worker service

[Service]
Type=forking
ExecStart=/usr/bin/laravel_queue start
ExecStop=/usr/bin/laravel_queue stop

[Install]
WantedBy=multi-user.target